# support/ldap-delegate.service
[Unit]
Description=ORM LDAP Delegate Service
After=multi-user.target
After=network.target
After=syslog.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'exec /opt/OpenResourceManager/LdapDelegate/ldap-delegate > >(/usr/bin/tee -a /var/log/orm/LdapDelegate/run.log) 2> >(/usr/bin/tee -a /var/log/orm/LdapDelegate/error.log >&2)'
Restart=always
RestartSec=2
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ldap-delegate

[Install]
WantedBy=multi-user.target